# autostart/lxqt-ptbatterysystemtray-autostart.desktop.in
[Desktop Entry]
Type=Application
Name=PtBatterySystemTray
TryExec=ptbatterysystemtray
Exec=ptbatterysystemtray
OnlyShowIn=LXQt;

#TRANSLATIONS_DIR=translations
