# autostart/lxqt-qxkb-autostart.desktop.in
[Desktop Entry]
Type=Application
Name=qxkb
TryExec=qxkb
Exec=qxkb
OnlyShowIn=LXQt;

#TRANSLATIONS_DIR=translations
